# sv/.log/run
#!/bin/sh

. /etc/runit/functions.bash

while read LINE; do
	logger -t "$(daemon_name)" "$LINE"
done
