# sv/slim/run
#!/bin/sh

exec slim < /dev/tty12 > /dev/tty12 2> /dev/tty12
