# sv/sshd/run
#!/bin/sh

exec /usr/host/bin/sshd -D 2>&1
