# sv/bluealsa/run
#!/bin/sh

. /etc/runit/functions.bash

needs bluetoothd
needs dbus-system

exec /usr/host/bin/bluealsa 2>&1
