# sv/dbus-system/run
#!/bin/sh

mkdir -p /run/dbus
rm -f /run/dbus.pid
exec dbus-daemon --system --nofork
