# sv/xinetd/run
#!/bin/sh

exec /usr/host/bin/xinetd -dontfork
