# sv/agetty-tty1/run
#!/bin/sh

exec chpst -P agetty -8 tty1
