# sv/bluetoothd/run
#!/bin/sh

needs dbus-system

exec /usr/host/libexec/bluetooth/bluetoothd -C -n -E
