# sv/cupsd/run
#!/bin/sh

. /etc/runit/functions.bash

exec /usr/host/bin/cupsd -f
