# sv/.dhcpcd/run
#!/bin/sh

. /etc/runit/functions.bash

IFNAME=$(param $0)

needs udevd
wants_down "dhcpd-$IFNAME"

exec dhcpcd -B $IFNAME
