# sv/.dhcpd/run
#!/bin/sh

. /etc/runit/functions.bash

IFNAME="$(param $0)"

needs udevd
wants_down "dhcpcd-$IFNAME"

exec dhcpd -f "$IFNAME"
