# sv/ntpd/run
#!/bin/sh

exec /usr/host/bin/ntpd -n -S "hwclock -w" 2>&1
