# sv/wm/finish
#!/bin/bash

killall sway || true
killall X || true

. ../.log/finish
